# user.timer
[Unit]
Description=Timer controlling the run of the users job

[Timer]
OnUnitActiveSec=10s
OnBootSec=10s

[Install]
WantedBy=timers.target
